# roles/prepare-servers/files/rc.local
#!/bin/bash
sudo ifup eth0
sudo ifup lo
exit 0

