# roles/prepare-oob-server/files/rc.local
#!/bin/bash
sudo bash /etc/network/if-pre-up.d/iptables
exit 0

